# src/dist/omero-ms-image-region.service
[Unit]
Description=OMERO image region microservice server
Requires=multi-user.target
After=multi-user.target
After=network.service

[Service]
Type=simple
Environment="JAVA_OPTS=-Dlogback.configurationFile=/opt/omero/OMERO.ms/omero-ms-image-region/current/conf/logback.xml"
WorkingDirectory=/opt/omero/OMERO.ms/omero-ms-image-region/current
ExecStart=/opt/omero/OMERO.ms/omero-ms-image-region/current/bin/omero-ms-image-region
User=omero
Group=omero
Restart=no
LimitNPROC=16384
LimitNOFILE=16384

[Install]
WantedBy=multi-user.target
